# cloudflared-api.service
[Unit]
Description=Cloudflare Tunnel for API Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/cloudflared tunnel --config /home/mark/ticketer-tracker/api-tunnel.yml run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target